# conf/service/provisioning/run
#!/bin/sh
exec 2>&1
exec su -c /home/bonitomail/scripts/activate-and-notify.sh bonitomail
